# contrib/sv.service
[Unit]
Description=SoundVibes speech-to-text daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/sv daemon start --vad continuous
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
