# contrib/sv.service
[Unit]
Description=SoundVibes speech-to-text daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/sv daemon start --vad continuous --vad-silence-ms 600
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
